# dashboard/config/prod/hub.conf
description "uWSGI server for OpenElections metadata hub"
 
start on runlevel [2345] # start on all runlevels.
stop on runlevel [!2345] # stop when shutting down.
 
respawn
exec /usr/bin/uwsgi -i /home/openelex/.virtualenvs/openelex-dash/src/dashboard/dashboard/config/prod/uwsgi.ini
